# toothless.conf
description "Toothless"

start on started docker
stop on stopping docker

respawn

script
    cd /opt/toothless
    export DISPLAY=:0
    exec node/bin/node --harmony toothless.js
end script
